# deploy/launchd/com.babygrow.feedback-reminder.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.babygrow.feedback-reminder</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/wangpengbo/Documents/Claude/Projects/baby_grow_helper/scripts/check_feedback_digest.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/wangpengbo/Documents/Claude/Projects/baby_grow_helper</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/babygrow-feedback-reminder.out</string>
  <key>StandardErrorPath</key>
  <string>/tmp/babygrow-feedback-reminder.err</string>
</dict>
</plist>
